# autoDeployment/autodeploy.timer
[Unit]
Desciption=Auto Deploy Timer

[Timer]
Unit=autodeploy.service
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
