# autoDeployment/autodeploy.service
[Unit]
Description=Auto deploy script

[Service]
ExecStart=/usr/bin/autodeploy.sh

[Install]
WantedBy=graphical.target
